# user/dunst.service
[Unit]
Description=Dunst notification manager

[Service]
ExecStart=/usr/bin/dunst
Restart=always
