# user/fusuma.service
[Unit]
Description=Fusuma gestures

[Service]
ExecStart=fusuma
