# user/mpd_volume.service
[Unit]
Description=Listener for mpd volume changes

[Service]
ExecStart=/usr/bin/ruby "${HOME}/scripts/mpd_volume.rb"
