# user/networkd-notify.service
[Unit]
Description=Networkd Notify

[Service]
ExecStart=/home/rein/projects/networkd-notify/networkd-notify.rb

[Install]
WantedBy=default.target
