# user/picom.service
[Unit]
Description=Picom

[Service]
ExecStart=/usr/bin/picom --config "${HOME}/.config/picom.conf"
