# user/picom.service
[Unit]
Description=Start picom

[Service]
ExecStart=/usr/bin/picom

[Install]
WantedBy=default.target
