# user/polybar_code.service
[Unit]
Description=Polybar code

[Service]
ExecStart=/usr/bin/polybar code
Restart=on-failure

[Install]
WantedBy=polybar.target
