# user/polybar@.service
[Unit]
Description=Polybar %i

[Service]
ExecStart=/usr/bin/polybar %i
Restart=on-failure
