# user/polybar.target
[Unit]
Description=Polybar target
After=graphical.target
