# user/polybar_main.service
[Unit]
Description=Polybar main

[Service]
ExecStart=/usr/bin/polybar main
Restart=on-failure

[Install]
WantedBy=polybar.target
