# user/synonym.timer
[Unit]
Description=Give a daily synonym of the day message

[Timer]
OnCalendar=*-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target
