# user/synonym.service
[Unit]
Description=Give synonym of the day message

[Service]
ExecStart=/bedrock/cross/bin/ruby %h/scripts/synonym.rb
