# user/synonym.service
[Unit]
Description=Give synonym of the day message

[Service]
Type=forking
ExecStart=/home/rein/scripts/synonym.rb
