# user/udiskie.service
[Unit]
Description=Start udiskie

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
