# user/udiskie.service
[Unit]
Description=Udiskie

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=graphical.target
