# user/unison@.service
[Unit]
Description=Sync directory %i with unison

[Service]
ExecStart=bash -c 'unison -sshargs "-i /home/rein/.ssh/xlab18" $HOME/%i ssh://x-lab18.local.boem/unison/%i -auto -batch -times'
Restart=no
